# quadlet/smarthome/rest980/php-nginx.container
[Container]
ContainerName=php-nginx
Environment=NGINX_WEBROOT=/app
Image=webhippie/php-nginx:latest
Volume=haconfig.volume/vacuum:/app:z

[Service]
Restart=always